# 04_db/atom_data/DETECT/detect_endpoint_wirekecr_001.yaml
# GR Atom: Windows Registry Key Creation Event
# Schema: atom_schema.yaml v5.1
# Status: enriched

# ==== Section 1: Identity ====
identity:
  id: "DETECT-ENDPOINT-WIREKECR-001"
  name: "Windows Registry Key Creation Event"
  normalized_name: "windows_registry_key_creation_event"
  aliases:
    - "Registry Key Create"
    - "RegCreateKeyEx Detection"
    - "New Registry Key Event"

# ==== Section 2: Classification ====
classification:
  is_infrastructure: false
  type: "pattern"
  abstraction_level: 1
  atom_tags:
    - EDR
    - SIEM
    - AUDIT
    - LOG
    - MONITOR
    - PERSIST
    - WINDOWS

# ==== Section 4: Scope (non-INFRA only) ====
scope:
  target_layers: ["L6", "L7"]
  target_zones: ["Z1", "Z2", "Z3"]
  target_components: []

# ==== Section 5: Definition ====
definition:
  what: "A detection pattern that captures the creation of new registry keys on a Windows endpoint. This is detected via Sysmon Event ID 12 (RegistryEvent - Object create and delete) when the event type indicates creation, and through Windows Security Event ID 4657 when registry object access auditing is enabled. The event records the full registry key path, the creating process (image path, PID), the user context, and the timestamp. Key creation under well-known persistence locations such as Run, RunOnce, Services, Winlogon, and AppInit_DLLs is particularly significant."
  why: "Registry key creation is a foundational operation for many adversary techniques including persistence establishment, configuration of malicious services, installation of browser helper objects, and setting up COM object hijacking. New keys under known persistence paths are strong indicators of compromise. Monitoring key creation also detects software installation activity, policy changes, and system configuration modifications that may affect security posture."
  how: "Detection monitors Sysmon Event ID 12 for registry key creation operations, filtering on EventType 'CreateKey'. Windows Security Event ID 4657 with access mask for KEY_CREATE_SUB_KEY provides additional coverage when SACLs are configured on monitored hives. EDR agents capture registry key creation with process context and enrichment. SIEM rules focus on creation events in high-value registry paths including HKLM\\SOFTWARE\\Microsoft\\Windows\\CurrentVersion\\Run, HKLM\\SYSTEM\\CurrentControlSet\\Services, and other persistence-related locations."
  core_concepts:
    - key: "persistence_key_creation"
      value: "Creation of registry keys in well-known autostart locations that enable code execution at boot or logon"
    - key: "process_attribution"
      value: "Identifying which process created the registry key to determine if the operation is legitimate or malicious"
    - key: "key_path_analysis"
      value: "Analyzing the full registry path to determine the purpose and security significance of the new key"
    - key: "temporal_correlation"
      value: "Correlating key creation with subsequent value writes to reconstruct the complete persistence setup sequence"

# ==== Section 6: Relations ====
relations:
  - type: "is_a"
    target: "DETECT-ENDPOINT-WINREGKEY-001"
    description: "Registry key creation is a specific type of registry key operation event"
  - type: "enables"
    target: "DETECT-ENDPOINT-WIREVASE-001"
    description: "Key creation enables subsequent registry value set operations within the new key"

# ==== Section 7: Security Profile (optional) ====
security_profile:
  mitre_mapping:
    technique_id: "T1547.001"
    tactic:
      - "persistence"
      - "privilege-escalation"
    sub_techniques: []
  defense_context:
    effectiveness: "high"

# ==== Section 8: Metadata ====
metadata:
  created: "2026-02-06"
  updated: "2026-02-10"
  version: "1.0"
  confidence: "high"
  trust_source: "primary"
  sources:
    - "Sysmon Documentation - Event ID 12 RegistryEvent"
    - "Microsoft Documentation - RegCreateKeyEx Function"
    - "MITRE ATT&CK T1547.001 - Boot or Logon Autostart Execution: Registry Run Keys"
    - "Windows Security Event ID 4657 Reference"
  search_keywords:
    - "registry key creation"
    - "RegCreateKeyEx"
    - "Sysmon Event ID 12"
    - "new registry key"
    - "registry persistence"
    - "Run key creation"
    - "registry autostart"
    - "registry key create"
    - "HKLM key creation"
    - "registry monitoring"
  embedding_text: "Windows registry key creation event detection pattern captures creation of new registry keys via Sysmon Event ID 12 and Windows Security Event ID 4657. Records the full key path, creating process, user context, and timestamp. Key creation under persistence locations such as Run, RunOnce, Services, Winlogon, and AppInit_DLLs is critical for detecting adversary persistence establishment. Maps to MITRE ATT&CK T1547.001 Registry Run Keys. Detection focuses on high-value registry paths with process attribution to distinguish legitimate software installation from malicious persistence setup and COM object hijacking."
